# crontab
0 3 * * * root /etc/cron.d/restart_docker_compose.sh
